# supervisord.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=rsyslogd

[program:sshd]
command=/usr/sbin/sshd -D

[program:dnet_ioserv]
command=dnet_ioserv -c /etc/elliptics/ioserv.json
